# nasa-apod.service
[Unit]
Description=NASA APOD Wallpaper Refresh

[Service]
ExecStart=PATH_TO_SCRIPT
